# profile.d/Kubernetes_helm_profile.ps1
#Set-Alias helmfile 'C:\Program Files (x86)\Helmfile\helmfile_windows_386.exe'